# build/mysql/profile.d/mysql_aliases.sh
#!/bin/bash

alias db="mysql -h mysql -u root --password=\$MYSQL_ROOT_PASSWORD -A \$MYSQL_DATABASE"
