# build/php/profile.d/composer_aliases.sh
#!/bin/bash

alias cda="composer dump-autoload"

