# build/php/profile.d/tempest_aliases.sh
#!/bin/bash

alias tempest="./tempest"
alias t="tempest"
