# build/php_old/profile.d/git_aliases.sh
#!/bin/bash

alias gf="echo git fetch --prune --all && git fetch --prune --all"
alias gpull="echo git pull && git pull"
alias gpush="echo git push && git push"
alias gdiff="echo git diff && git diff"
alias ga="echo git add . --all && git add . --all"
alias gs="echo git status && git status"
alias gb="echo git branch && git branch"
alias gc="echo git checkout XXX && git checkout"
alias gcm="echo git checkout master && git checkout master"
alias gcp="echo git checkout production && git checkout production"
alias gmm="echo git merge master && git merge master"
