# build/php_claude/profile.d/tempest_aliases.sh
#!/bin/bash

alias tempest="./tempest"
alias tmp="tempest"
